# test_monitor.service
[Unit]
Description=Monitoring test process

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test_monitor.sh
